# dot_config/systemd/user/safeeyes-apply-schedule.service
[Unit]
Description=Apply safeeyes work-hours schedule

[Service]
Type=oneshot
ExecStart=%h/.local/bin/safeeyes-apply-schedule
